# etc/supervisor/conf.d/ffmpeg.speaker.conf
[program:ffmpeg.speaker]
priority=40
command=/composer/ffmpeg.speaker.sh
autostart=false
autorestart=true
startsecs=10
stopsignal=QUIT
redirect_stderr=true
